# deploy/paul/qgis_server_9990.conf
description "FastCGI QGIS Server for helloworld project"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 3 5

setuid paul
setgid paul

script
    exec /usr/bin/spawn-fcgi \
        -n \
        -d /home/paul/wk/cccs/core/qgis/projects/helloworld \
        -p 9990 \
        -- \
        /usr/lib/cgi-bin/qgis_mapserv.fcgi
end script