# _deploy/staging/upstart.conf
start on started postgresql
stop on stopped postgresql

respawn
respawn limit 3 5

setuid webcore
setgid webcore

exec /home/webcore/staging/_deploy/staging/service_demon.sh