# _deploy/cristi/upstart.conf
start on started rc
stop on stopped rc

respawn
respawn limit 3 5

setuid cristi
setgid cristi

exec /home/cristi/cccs_web/_deploy/cristi/service_demon.sh