# deploy/production/upstart.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 3 5

setuid cccs
setgid cccs

exec /home/cccs/production/deploy/production/service_demon.sh