# _deploy/production_upstart.conf
start on started rc
stop on stopped rc

respawn
respawn limit 3 5

setuid webcore
setgid webcore

exec /home/webcore/production/_deploy/production_service_demon.sh