# _deploy/production_upstart.conf
start on started postgresql
stop on stopped postgresql

respawn
respawn limit 3 5

setuid webcore
setgid webcore

exec /home/webcore/core/_deploy/production_service_demon.sh