# deploy/staging/upstart.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 3 5

setuid webcore
setgid webcore

exec /home/webcore/staging/deploy/staging/service_demon.sh