# deploy/cristi/upstart.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 3 5

setuid cristi
setgid cristi

exec /home/cristi/cccs_web/deploy/cristi/service_demon.sh