# deploy/test/upstart.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 3 5

setuid cccs
setgid cccs

exec /home/cccs/test/deploy/test/service_demon.sh